# .docker/s6-rc.d/cp2k-code-setup/up
#!/command/execlineb -S0

with-contenv

foreground { s6-echo "Setting up CP2K code" }

/etc/init/add-codes.sh
